# DEBIAN/prerm
#!/bin/sh -e
/etc/init.d/mmonit stop || true 
update-rc.d -f mmonit remove
